# supervisord.conf
[supervisord]
nodaemon=true

[program:kibana]
command=/opt/kibana-5.3.2-linux-x86_64/bin/kibana
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
autorestart=true
